# docker/root/etc/services.d/xdsgenerator/run
#!/usr/bin/execlineb -P

with-contenv
s6-envuidgid opentele
fdmove -c 2 1
s6-envdir /etc/services.d/xdsgenerator/envdir
cd /app
java -jar xds-generator.jar
