# docker/root/etc/services.d/xdsgenerator/run
#!/command/execlineb -P

with-contenv
s6-envuidgid nobody
fdmove -c 2 1
s6-envdir /etc/services.d/xdsgenerator/envdir
cd /app
java -jar xds-generator.jar
